# opt/local/etc/profile.d/00-banner.sh
#!/bin/sh

[ -f /etc/banner ] && cat /etc/banner

[ -e /tmp/.failsafe ] && cat /etc/banner.failsafe

