# opt/local/etc/profile.d/20-mksh-config.sh

###################################
if [ ! -z "$KSH_VERSION" ] ; then #
###################################

set -o utf8-mode

bind '^L=clear-screen'


####
fi #
####

