# opt/local/etc/profile.d/opt.sh.off
_OPTLOCAL_ENV='/opt/local/lib/common.sh'
[ -f "$_OPTLOCAL_ENV" ] && \
  alias optlocalenv=". $_OPTLOCAL_ENV"
unset _OPTLOCAL_ENV
