# opt/local/etc/rc.local
#!/bin/sh

[ -f /opt/local/etc/rc.bootstrap ] && {
    . /opt/local/etc/rc.bootstrap
}

exit 0
